# android-app/agent/com.minerva.agent-android.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.minerva.agent-android</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jairo/minerva/android-app/agent/run-agent-mac.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/minerva-agent-android.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/minerva-agent-android.log</string>
</dict>
</plist>
